# pharma_gyan_proj/server_config/services/pharma_gyan_web_service.conf
[unix_http_server]
file=/dev/shm/supervisor.sock   ; (the path to the socket file)

[supervisord]
logfile=/opt/logs/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=1024MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=3           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/dev/shm/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=100                  ; (min. avail startup file descriptors;default 1024)
minprocs=500                 ; (min. avail process descriptors;default 200)
[rpcinterface:supervisor]
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///dev/shm/supervisor.sock ; use a unix:// URL  for a unix socket

[supervisord]
nodaemon=true

[program:pharma_gyan]
environment="NEW_RELIC_ENVIRONMENT=web"
environment="NEW_RELIC_CONFIG_FILE=/etc/newrelic/newrelic_pharma_gyan.ini"
environment=
         C_FORCE_ROOT=true
# command=bash -c "/usr/local/python3.8/bin/uwsgi --ini /etc/pharma_gyan_uwsgi.ini"
command=bash -c "export NEW_RELIC_ENVIRONMENT=web && export NEW_RELIC_CONFIG_FILE=/etc/newrelic/newrelic_pharma_gyan.ini && /usr/local/python3.8/bin/newrelic-admin run-program /usr/local/python3.8/bin/uwsgi --ini /etc/pharma_gyan_uwsgi.ini"
autostart=true
autorestart=unexpected
logfile_maxbytes=0        ; max # logfile bytes b4 rotation (default 50MB)
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stdout_logfile_backups=3
stderr_logfile=/dev/fd/1
stderr_logfile_maxbytes=0
exitcodes=1

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/1
stderr_logfile_maxbytes=0


